# scripts/piaware.service
# piaware uploader service for systemd
# install in /etc/systemd/system

[Unit]
Description=FlightAware ADS-B uploader
Documentation=https://flightaware.com/adsb/piaware/
Wants=network-online.target
After=dump1090-fa.service network-online.target time-sync.target

[Service]
ExecStart=/usr/bin/piaware -p /var/run/piaware.pid -plainlog
Type=simple
Restart=on-failure
RestartSec=30

[Install]
WantedBy=default.target
